# lib/systemd/system/solydxk-root-gtk3.service
[Unit]
Description=Root gtk3

[Service]
ExecStart=/usr/lib/solydxk/scripts/root-gtk3

[Install]
WantedBy=multi-user.target 
